# configs/poolserver.conf
description "poolserver"

start on filesystem
stop on runlevel [!2345]
respawn
respawn limit 1000 600 # 10 times in 600 seconds
console log

setuid poolserver
exec /usr/bin/PoolServer
